# .mine/etc/profile.d/aws_login.sh
aws_login() {
	echo "from bash"
	powershell awssec -u fabien_hiegel
}
