# .mine/etc/profile.d/git-layout.sh
#
# CONFIGURATION DU PROMPT GIT
#
source "${HOME_DIR}/.gitsettings/layout"
