# etc/profile.d/git-layout.sh
#
# CONFIGURATION DU PROMPT GIT
#
source "${MINE_DIR}/etc/git/layout"
