# etc/profile.d/opera.sh
alias opera='"C:\Users\Fabien_Hiegel\AppData\Local\Programs\Opera\launcher.exe" "--user-data-dir=C:\Users\Fabien_Hiegel\fhiegel\Opera2"'
