# etc/profile.d/ubuntu.sh
alias ubuntu='cd /c/Users/Fabien_Hiegel/AppData/Local/Packages/CanonicalGroupLimited.Ubuntu18.04onWindows_79rhkp1fndgsc/LocalState/rootfs'
